# crontab
* */10 * * * * * python UBA_liveparser.py 
